# docker/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
user=docker
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=php-fpm -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=0

#[program:schedule]
#process_name=%(program_name)s
#command=php /var/www/artisan schedule:work
#autostart=true
#autorestart=true
#redirect_stderr=true
#stdout_logfile=/var/www/storage/logs/schedule.log
#stopwaitsecs=3600

[program:websockets]
command=php /var/www/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/websokets.log
user=docker

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf
